# crontab
0 3 * * * bundle exec rake push_line:push_line_message_notify_date
0 7 * * * bundle exec rake push_line:push_line_message_notify_date