# crontab
0 22 * * * bundle exec rake push_line:push_line_message_notify_date